# root/etc/services.d/adguard/run
#!/bin/bash

exec /usr/sbin/adguard -c /config/adguard/conf/AdGuardHome.yaml -w /config/adguard/data